# sup-docker.conf
[supervisord]
nodaemon=true
user=root

logfile_maxbytes=5MB ; (max main logfile bytes b4 rotation;default 50MB)

logfile_backups=5 ; (num of main logfile rotation backups;default 10)
loglevel=INFO ; (log level;default info; others: debug,warn,trace)
minfds=10240 ; (min. avail startup file descriptors;default 1024)
minprocs=1000 ; (min. avail process descriptors;default 200)

# [inet_http_server]         
# port=127.0.0.1:9001

# [rpcinterface:supervisor]
# supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# [supervisorctl]
# serverurl=http://127.0.0.1:9001

#logfile_maxbytes=5MB ; (max main logfile bytes b4 rotation;default 50MB)

#logfile_backups=5 ; (num of main logfile rotation backups;default 10)
#loglevel=warn ; (log level;default info; others: debug,warn,trace)
#pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
#nodaemon=false ; (start in foreground if true;default false)
#minfds=10240 ; (min. avail startup file descriptors;default 1024)
#minprocs=1000 ; (min. avail process descriptors;default 200)


[program:createJson]
priority=1
numprocs=1
# command=/home/ubuntu/boTrading/bin/python /home/ubuntu/boTrading/backend/createJson.py
command=/usr/local/bin/python3 /usr/src/backend/createJson.py
autostart=true
autorestart=true
# startsecs=1
user=root
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true
#environment=PYTHONPATH=/usr/local/lib/python38.zip:/usr/local/lib/python3.8:/usr/local/lib/python3.8/site-packages