# supervisord.conf
[supervisord]

logfile_maxbytes=5MB ; (max main logfile bytes b4 rotation;default 50MB)

logfile_backups=5 ; (num of main logfile rotation backups;default 10)
loglevel=warn ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false ; (start in foreground if true;default false)
minfds=10240 ; (min. avail startup file descriptors;default 1024)
minprocs=1000 ; (min. avail process descriptors;default 200)


#logfile_maxbytes=5MB ; (max main logfile bytes b4 rotation;default 50MB)

#logfile_backups=5 ; (num of main logfile rotation backups;default 10)
#loglevel=warn ; (log level;default info; others: debug,warn,trace)
#pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
#nodaemon=false ; (start in foreground if true;default false)
#minfds=10240 ; (min. avail startup file descriptors;default 1024)
#minprocs=1000 ; (min. avail process descriptors;default 200)


[program:createJson]
priority=1
numprocs=1
command=/home/ubuntu/boTrading/bin/python /home/ubuntu/boTrading/backend/createJson.py
autostart=true
autorestart=true

[program:saveOdds]
priority=1
numprocs=1
command=/home/ubuntu/boTrading/bin/python /home/ubuntu/boTrading/backend/saveOdds.py
autostart=true
autorestart=true